# sources/debian/sambaedu-web-common.cron.d
0 * * * * www-admin -c "kinit -k -t /var/www-admin/www-sambaedu.keytab www-sambaedu"
