# sources/debian/sambaedu-web-common.cron.d
0 * * * *	www-admin	"kinit -k -t /etc/sambaedu/www-sambaedu.keytab www-sambaedu"
@reboot	www-admin	"kinit -k -t /etc/sambaedu/www-sambaedu.keytab www-sambaedu"
