# sources/debian/sambaedu-web-common.cron.d
50 * * * * www-admin "/usr/share/sambaedu/sbin/renew_ticket.sh"
@reboot www-admin "/usr/share/sambaedu/sbin/renew_ticket.sh"
